# scripts/systemd/market-intelligence.service
[Unit]
Description=NL Bots Market Intelligence Update
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/nl-international-ai-bots
Environment="PYTHONPATH=/root/nl-international-ai-bots"
ExecStart=/root/nl-international-ai-bots/venv/bin/python /root/nl-international-ai-bots/scripts/update_market_intelligence.py

# Логирование
StandardOutput=journal
StandardError=journal

# Перезапуск при ошибке
Restart=on-failure
RestartSec=5min

[Install]
WantedBy=multi-user.target
