# deploy/partner-panel/partner-panel.service
[Unit]
Description=NL Partner Panel Backend (FastAPI)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/apexflow

# Переменные окружения
Environment="PATH=/root/apexflow/venv/bin"
Environment="PARTNER_PANEL_API_PORT=8001"
Environment="PARTNER_PANEL_DATABASE_URL=postgresql+asyncpg://postgres:UB8TG6@@IUYDGC@localhost:5432/nl_international"
Environment="PARTNER_PANEL_SECRET_KEY=your-super-secret-key-change-this"
Environment="PARTNER_PANEL_BOT_TOKEN=YOUR_BOT_TOKEN_HERE"

# Запуск uvicorn
ExecStart=/root/apexflow/venv/bin/uvicorn partner_panel.backend.main:app --host 127.0.0.1 --port 8001 --workers 2

# Перезапуск при падении
Restart=always
RestartSec=5

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=partner-panel

[Install]
WantedBy=multi-user.target
