# scripts/apexflow-traffic.service
[Unit]
Description=APEXFLOW Traffic Engine - NL International
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/apexflow
ExecStart=/root/apexflow/venv/bin/python run_auto_comments.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/root/apexflow/.env

# Limits
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
